# systemd/trovix-capi-queue.target
[Unit]
Description=Trovix Queue Target
After=redis.service

[Install]
WantedBy=multi-user.target
